# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true

[program:dotnet]
command=/usr/bin/dotnet /app/ARM-Docker-Api-Deploy.dll
autostart=true
autorestart=true
